# app-misc/wmr-oregon/files/wmrd
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

extra_started_commands="reload reinit"
command="/usr/sbin/wmr-oregon"
VMRD_PID="/var/run/wmrd.pid"
command_args="-d -l ${VMRD_PID}"
start_stop_daemon_args="--quiet -p ${VMRD_PID} -b"
description="Daemon for Advanced Configuration and Power Interface"

description_reload="Reloads the configuration."
description_reinit="Re-inits related USB devices and reloads the configuration."

depend() {
    need localmount
    use logger
}

reload() {
    ebegin "Reloading wmrd configuration"
    start-stop-daemon --exec $command --signal HUP
    eend $?
}

reinit() {
    ebegin "Reiniting related USB devices and reloading wmrd configuration"
    start-stop-daemon --exec $command --signal USR1
    eend $?
}

